# deploy/dmp-packing.service
[Unit]
Description=DMP Container Packing Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/dmp-packing-service
Environment="PATH=/var/www/dmp-packing-service/venv/bin"
ExecStart=/var/www/dmp-packing-service/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dmp-packing

[Install]
WantedBy=multi-user.target
